# packaging/vsync-locker/macos/com.field.vsync-locker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.field.vsync-locker</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/python3</string>
    <string>__PREFIX__/lib/field-vsync-locker.py</string>
    <string>guard</string>
    <string>--quiet</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>NEXUS_INSTALL_ROOT</key>
    <string>__PREFIX__</string>
    <key>NEXUS_STATE_DIR</key>
    <string>__PREFIX__/.nexus-state</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
</dict>
</plist>